# deploy/systemd/rook-servicechannel-gateway.service
[Unit]
Description=RooK Servicechannel Gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=rook-gateway
Group=rook-gateway
WorkingDirectory=/opt/rook-servicechannel-gateway
EnvironmentFile=/etc/rook-servicechannel-gateway/gateway.env
ExecStart=/usr/local/bin/rook-servicechannel-gateway
Restart=on-failure
RestartSec=5s
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadOnlyPaths=/
ReadWritePaths=/var/log/rook-servicechannel-gateway

[Install]
WantedBy=multi-user.target
